# supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=jupyter notebook --ip 0.0.0.0 --allow-root --port 8888
stderr_logfile = /var/log/supervisord/jupyter-stderr.log
stdout_logfile = /var/log/supervisord/jupyter-stdout.log
startsecs = 0

[program:flask]
command=flask run --port 5000
startsecs = 0
stderr_logfile = /var/log/supervisord/flask-stderr.log
stdout_logfile = /var/log/supervisord/flask-stdout.log
